# templates/cron-daemon/systemd.service
# ApexRouter-CLI Cron Daemon — systemd user unit (Linux)
#
# Installation:
#   1. Copy this file to ~/.config/systemd/user/apexrouter-cron-daemon.service
#   2. Edit ExecStart below to point to the apexrouter-cli binary.
#   3. Optionally override APEXROUTER_CLI_HOME if it differs from ~/.apexrouter.
#   4. Enable + start:
#        systemctl --user daemon-reload
#        systemctl --user enable --now apexrouter-cron-daemon.service
#   5. Verify:
#        systemctl --user status apexrouter-cron-daemon.service
#
# Uninstall:
#   systemctl --user disable --now apexrouter-cron-daemon.service
#   rm ~/.config/systemd/user/apexrouter-cron-daemon.service
#   systemctl --user daemon-reload
#
# Logs:
#   journalctl --user -u apexrouter-cron-daemon.service -f
#   # or from the log file:
#   tail -f ~/.apexrouter/cron-daemon.log

[Unit]
Description=ApexRouter-CLI Cron Daemon
Documentation=https://github.com/APEX-AI-LABS-LLP/apexrouter-cli
After=network.target

[Service]
Type=simple

# Replace with the absolute path to the apexrouter-cli binary.
# Example: ExecStart=/usr/local/bin/apexrouter-cli cron daemon
ExecStart=/usr/local/bin/apexrouter-cli cron daemon

# Run the daemon body directly (systemd is the process manager).
Environment=APEXROUTER_CLI_CRON_DAEMON_CHILD=1

# Override if your ApexRouter home is not ~/.apexrouter.
# Environment=APEXROUTER_CLI_HOME=/home/youruser/.apexrouter

# Redirect logs to the daemon log file. systemd also captures
# stdout/stderr via journald by default.
StandardOutput=append:%h/.apexrouter/cron-daemon.log
StandardError=append:%h/.apexrouter/cron-daemon.log

# Restart on unexpected exit. Do not restart on SIGTERM (clean shutdown).
Restart=on-failure
RestartSec=10s

# Allow the unit to start at login via `loginctl enable-linger`.
KillMode=process
TimeoutStopSec=10

[Install]
WantedBy=default.target
